# software/FBRoot/etc/systemd/system/checknetwork.service
[Unit]
Description=Check network, restart dhcpcd if not working

[Service]
Type=oneshot
ExecStart=/usr/local/bin/checknetwork

