# software/FBRoot/etc/systemd/system/ntpd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ntpd -g

